# devops/ansible/girder.conf
# tangelo - tangelo job file

description "tangelo"
author "Jeff Baumes <jeff.baumes@kitware.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Essentially lets upstart know the process will detach itself to the background
# expect daemon

env PYTHONPATH=/vagrant/girder

# Start the process
script
  exec python -m girder
end script
