# devops/ansible/romanesco.conf
# romanesco - romanesco job file

description "romanesco"
author "Jeff Baumes <jeff.baumes@kitware.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

env PYTHONPATH=/home/vagrant/vtk-precise64/lib/python2.7/site-packages:/home/vagrant/vtk-precise64/lib:/vagrant/romanesco
env LD_LIBRARY_PATH=/home/vagrant/vtk-precise64/lib

setuid vagrant
setgid vagrant

# Start the process
script
  exec python -m romanesco
end script
